# config/systemd/user/battery-monitor.timer
[Unit]
Description=Run battery monitor every minute

[Timer]
OnBootSec=60
OnUnitActiveSec=60

[Install]
WantedBy=timers.target